# platform/broadcom/sonic-platform-modules-cel/debian/platform-modules-silverstone-dp.postinst
depmod -a
systemctl enable platform-modules-silverstone-dp.service
systemctl start platform-modules-silverstone-dp.service