# platform/broadcom/sonic-platform-modules-cel/debian/platform-modules-mystone.postinst
depmod -a

# Enable bmc vlan_ip service
systemctl enable bmc_vlan.service

systemctl enable platform-modules-mystone.service
systemctl start platform-modules-mystone.service

# Enable platform-sync_bmc-timer
# systemctl enable sync_bmc.timer
# systemctl start sync_bmc.timer

# Enable heartbeat timer
systemctl enable bmc_wdt.service
systemctl start bmc_wdt.service
