# platform/broadcom/sonic-platform-modules-cel/debian/platform-modules-questone2bd.postinst
depmod -a

systemctl enable platform-modules-questone2bd.service
systemctl start platform-modules-questone2bd.service
