# platform/broadcom/sonic-platform-modules-cel/debian/platform-modules-fishbone48.postinst
depmod -a

# Enable bmc vlan_ip service
systemctl enable bmc_vlan.service

update-rc.d platform-modules-fishbone48 defaults
service platform-modules-fishbone48 start

# Enable platform-sync_bmc-timer
# systemctl enable sync_bmc.timer
# systemctl start sync_bmc.timer