# platform/broadcom/sonic-platform-modules-cel/debian/platform-modules-questone2bd.postinst
depmod -a

systemctl enable platform-modules-questone2bd.service
systemctl start platform-modules-questone2bd.service

systemctl enable fanctrl.service
systemctl start fanctrl.service

/usr/local/bin/platform_api_mgnt.sh install
