# platform/broadcom/sonic-platform-modules-cel/debian/platform-modules-silverstone2.postinst
depmod -a
systemctl enable platform-modules-silverstone2.service
systemctl start platform-modules-silverstone2.service
