# platform/broadcom/sonic-platform-modules-cel/debian/platform-modules-mystone/etc/systemd/system/bmc_vlan.service
[Unit]
Description="Service for add vlan for rsyslog"
After=network.target
Before=rsyslog-config.service
Before=ntp.service

[Service]
Type=forking
ExecStart=/bin/sh /usr/local/etc/bmc_vlan.sh

[Install]
WantedBy=rsyslog.service