# platform/broadcom/sonic-platform-modules-cel/tools/sync_bmc/sync_bmc.timer
[Unit]
Description=Timer send sensor data value to BMC.

[Timer]
OnUnitActiveSec=5s
OnBootSec=5s
AccuracySec=1us

[Install]
WantedBy=timers.target